# multi-app1/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:redis-server]
command=/usr/bin/redis-server /etc/redis/redis.conf

[program:node]
command=node /usr/local/src/test.js
